# archiso/airootfs/etc/systemd/system/harness-firstboot.service
[Unit]
Description=HarnessOS First Boot Welcome
After=multi-user.target
ConditionPathExists=!/var/lib/harness-firstboot-done

[Service]
Type=oneshot
ExecStart=/usr/local/bin/harness-welcome
ExecStartPost=/usr/bin/touch /var/lib/harness-firstboot-done
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
